# printer-service/debian/postinst
#!/bin/sh
set -e

. /usr/share/debconf/confmodule

case "$1" in
    configure)
        db_get printer-service/ws-url
        WS_URL="$RET"

        db_get printer-service/printer-mode
        MODE="$RET"

        ARGS="--url \"${WS_URL}\""

        case "$MODE" in
            mock)
                ARGS="$ARGS --mock"
                ;;
            network)
                db_get printer-service/network-ip
                NET_IP="$RET"
                db_get printer-service/network-port
                NET_PORT="$RET"
                ARGS="$ARGS --ip \"${NET_IP}\" --port ${NET_PORT}"
                ;;
            usb)
                # No extra args needed for USB mode
                ;;
        esac

        cat > /etc/default/printer-service <<EOF
# Generated by dpkg configure - edit with: dpkg-reconfigure printer-service
PRINTER_SERVICE_ARGS=${ARGS}
EOF

        deb-systemd-helper enable printer-service.service || true
        deb-systemd-invoke restart printer-service.service || true
        ;;
esac

#DEBHELPER#
